# systemd/trade1-analysis.service
[Unit]
Description=trade-1 analysis and optimization run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/trade-1
EnvironmentFile=/opt/trade-1/.env
ExecStart=/opt/trade-1/scripts/run_analysis.sh
MemoryMax=800M
MemorySwapMax=4G
CPUQuota=90%
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/trade-1/data /opt/trade-1/config /opt/trade-1/logs
