# deploy/trade-learning-sync.timer
[Unit]
Description=Run trade-1 learning database sync every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=trade-learning-sync.service

[Install]
WantedBy=timers.target
