# deploy/trade-backup.timer
[Unit]
Description=Daily trade-1 backup

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
