# deploy/trade-learning-review-daily.timer
[Unit]
Description=Run trade-1 daily learning review

[Timer]
OnCalendar=*-*-* 00:10:00 UTC
AccuracySec=5min
Persistent=true
Unit=trade-learning-review-daily.service

[Install]
WantedBy=timers.target
