# deploy/trade-learning-review-monthly.timer
[Unit]
Description=Run trade-1 monthly learning review

[Timer]
OnCalendar=*-*-01 00:30:00 UTC
AccuracySec=15min
Persistent=true
Unit=trade-learning-review-monthly.service

[Install]
WantedBy=timers.target
